# .config/profile.d/sd.sh
export SD_ROOT=$HOME/bin/sd
export SD_CAT=bat
